# bindplane.service
/etc/systemd/system # more bindplane.service
################################
BindPlane systemd service      #
# Author: Olu Lawrence         #
# Date: Nov 12, 2025           #
################################

[Unit]
Description=BindPlane Server with Encrypted Configuration
Documentation=https://docs.bindplane.com
After=network.target

[Service]
Type=simple
User=root
Group=root

# Step 1: Deploy decrypted config before starting
ExecStartPre=/usr/local/bin/bindplane-config-manager.sh prepare

# Step 2: Start BindPlane
ExecStart=/usr/local/bin/bindplane serve --config /etc/bindplane/config.yaml

# Step 3: Wait for startup, then delete config file (security measure)
# Run as separate process so it doesn't block the service
ExecStartPost=/bin/bash -c 'sleep 15 && /usr/local/bin/bindplane-config-manager.sh cleanup &'

# Cleanup config on stop
ExecStopPost=/usr/local/bin/bindplane-config-manager.sh cleanup

# Restart configuration
#Restart=on-failure
#RestartSec=10s
#StartLimitInterval=5min
#StartLimitBurst=3

# Environment
Environment="BINDPLANE_CONFIG_PATH=/etc/bindplane"
WorkingDirectory=/var/lib/bindplane

# Security settings
#NoNewPrivileges=true
#PrivateTmp=true
#ProtectSystem=strict
#ProtectHome=true
#ReadWritePaths=/var/lib/bindplane /var/log/bindplane

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bindplane

#ReadWritePaths=/var/lib/bindplane /var/log/bindplane

[Install]
WantedBy=multi-user.target
